# polymarket-weather-bot.service
[Unit]
Description=Polymarket Weather Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-weather-bot
Environment=BOT_PORT=3000
Environment=BOT_POLL_SECONDS=300
Environment=BOT_DB_PATH=/home/ubuntu/polymarket-weather-bot/bot.db
ExecStart=/usr/bin/python3 /home/ubuntu/polymarket-weather-bot/run_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
